# Examples/wiredbot.service
# systemd service unit for WiredChatBot
# Install: cp wiredbot.service /etc/systemd/system/
#          systemctl daemon-reload
#          systemctl enable --now wiredbot

[Unit]
Description=WiredChatBot — AI chatbot for Wired 3
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/var/run/wiredbot.pid
ExecStart=/usr/local/bin/WiredChatBot run \
    --config /etc/wiredbot/wiredbot.json \
    --spec   /usr/share/wiredbot/wired.xml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
User=wiredbot
Group=wiredbot

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/run /var/log/wiredbot

[Install]
WantedBy=multi-user.target
